# package/freebsd/usr/ports/security/p5-openxpki-deployment/files/openxpki.in
#!/bin/sh
# $FreeBSD: ports/security/p5-openxpki-deployment/files/openxpki.sh.in,v 1.2 2007/02/03 09:49:52 miwi Exp $

# PROVIDE: openxpki
# REQUIRE: LOGIN
# REQUIRE: postgresql
# REQUIRE: mysql
# KEYWORD: shutdown

# Define this variable in the file /etc/rc.conf:
# openxpki_enable="YES"
# to ensure that openxpki starts at boot time.
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#

. /etc/rc.subr

openxpki_enable="${openxpki_enable-NO}"
openxpki_flags=""
openxpki_pidfile="%%PREFIX%%/var/openxpki/openxpki.pid"

case "$1" in 
    *start) command_args="start";;
esac

name="openxpki"
rcvar=`set_rcvar`
command="%%PREFIX%%/bin/openxpkictl"
procname="openxpkid"

load_rc_config $name

pidfile="${openxpki_pidfile}"

run_rc_command "$1"
